# example/systemd/rusticlone.timer
[Unit]
Description=Rusticlone Backup - timer

[Timer]
# every day at midnight
OnCalendar=*-*-* 00:00:00

[Install]
WantedBy=timers.target
